# deploy/config/supervisord.conf.template
[supervisord]
user=root
logfile=/var/log/supervisord.log
pidfile=/var/log/supervisord.pid
nodaemon=true

[program:xvfb]
command=/bin/sh -c "rm -f /tmp/.X1-lock /tmp/.X11-unix/X1; mkdir -p /tmp/.X11-unix; exec /usr/bin/Xvfb :1 -screen 0 1280x800x24"
autostart=true
autorestart=true
priority=10
stderr_logfile=/var/log/xvfb.err.log
stdout_logfile=/var/log/xvfb.out.log
environment=DISPLAY=":1"

[program:app]
command=copaw app --host 0.0.0.0 --port ${COPAW_PORT}
autostart=true
autorestart=true
priority=30
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log
environment=DISPLAY=":1",PLAYWRIGHT_CHROMIUM_EXECUTABLE_PATH="/usr/bin/chromium",COPAW_RUNNING_IN_CONTAINER="1"
